# templates/default/verdaccio.service.erb
# Automatically Generated by Chef

[Unit]
Description=Verdaccio NPM cache & private repo server

[Service]
ExecStart=<%= @executable %> -c <%= @configfile %>
Type=simple
User=verdaccio
Group=verdaccio

[Install]
WantedBy=multi-user.target
